# files/target/etc/cron.weekly/puppet-etcd-defrag
#!/bin/bash

RANDOM_SLEEP=$[( $RANDOM % 3600 )]
echo "Sleeping ${RANDOM_SLEEP} seconds...."

ETCDCTL_API=3 etcdctl defrag
